# portal/portal.service
[Unit]
Description=ReSuture FluoroSim customer portal (gunicorn)
After=network-online.target
Wants=network-online.target

[Service]
# Deploy path: adjust if the repo is cloned elsewhere.
WorkingDirectory=/opt/fluoro_simulator
# All secrets/config live here (root:root, chmod 600) — see portal/README.md.
EnvironmentFile=/etc/resuture-portal.env
ExecStart=/opt/fluoro_simulator/portal/venv/bin/gunicorn -w 2 -b 127.0.0.1:8080 portal.app:create_app()
Restart=on-failure
RestartSec=3
# Run as an unprivileged user (create with: useradd -r -s /usr/sbin/nologin portal).
User=portal
Group=portal

[Install]
WantedBy=multi-user.target
